# virsh-api-server.service
[Unit]
Description=Virsh-Api-Server

[Service]
ExecStart=/Virsh-Api-Server/start.sh

[Install]
WantedBy=multi-user.target